# rootfs/etc/services.d/deluge-web/run
#!/usr/bin/execlineb -P

/usr/bin/deluge-web -c /config
