# rootfs/etc/s6-overlay/s6-rc.d/deluge-web/run
#!/command/execlineb -P

s6-setuidgid deluge

/usr/bin/deluge-web
    --config /config \
    --interface 0.0.0.0 \
    --loglevel info \
    --logfile /config/deluge-web.log \
    --pidfile /tmp/deluge-web.pid \
    --do-not-daemonize
