# rootfs/etc/s6-overlay/s6-rc.d/deluge-init/up
#!/command/execlineb -P

/bin/chown 1000:1000 /config -R
