# rootfs/etc/s6-overlay/s6-rc.d/deluged/run
#!/command/execlineb -P

s6-setuidgid deluge

/usr/bin/deluged \
    --config /config \
    --ui-interface 0.0.0.0 \
    --interface 0.0.0.0 \
    --loglevel info \
    --logfile /config/deluged.log \
    --pidfile /run/deluge.pid \
    --do-not-daemonize
