# rootfs/etc/s6-overlay/s6-rc.d/deluged/run
#!/usr/bin/execlineb -P

/usr/bin/deluged \
    --config /config \
    --ui-interface 0.0.0.0 \
    --interface 0.0.0.0 \
    --loglevel info \
    --logfile /dev/stdout \
    --user deluge \
    --group deluge \
    --do-not-daemonize
